# scripts/com.wenest.blog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Wenest blog — automatic publisher schedule (macOS launchd).
  Publishes one article on Monday, Wednesday and Friday at 10:00 LOCAL time.
  (Local = your Mac's timezone. The exact hour does not matter for SEO; it
  only needs the Mac powered on around then.)

  ACTIVATE (run once in Terminal, after you're happy with a dry-run):
    cp "/Users/cristianvwz/BusinessAI/WeNest/wenest-blog/scripts/com.wenest.blog.plist" ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.wenest.blog.plist

  DEACTIVATE:
    launchctl unload ~/Library/LaunchAgents/com.wenest.blog.plist

  Logs: /tmp/wenest-blog.log (output) and /tmp/wenest-blog.err (errors).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wenest.blog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/cristianvwz/BusinessAI/WeNest/wenest-blog/scripts/run-auto-publish.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>10</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>10</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>10</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/wenest-blog.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/wenest-blog.err</string>
</dict>
</plist>
